# units/platformd-secretd.service.in
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=platformd Secret Service provider
Documentation=man:platformd-secretd.service(8)

[Service]
Type=notify
ExecStart=@LIBEXECDIR@/platformd-secretd
Restart=on-failure
LimitCORE=0
MemorySwapMax=0

# Encryption at rest (optional): a systemd-creds-encrypted 32-byte vault key,
# delivered as a credential the daemon reads at $CREDENTIALS_DIRECTORY/vault-key.
# With no key configured the store stays in the clear (rely on an encrypted home).
# Create once, then uncomment:
#   head -c 32 /dev/urandom > k && systemd-creds encrypt k vault-key.cred && rm k
#LoadCredentialEncrypted=vault-key:%h/.config/platformd-secretd/vault-key.cred

# systemd-native sandboxing (kept moderate; nothing here blocks writing the
# user's own home, where the store lives)
NoNewPrivileges=yes
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=graphical-session.target
